# cron/cron.hourly/run_hourly
#!/bin/sh

/usr/local/bin/python3 /data/InStock/instock/job/basic_data_daily_job.py
